# 12h.timer
[Timer]
OnBootSec=0
OnUnitActiveSec=12h 15min

[Install]
WantedBy=timers.target
