# ptc.timer
[Timer]
OnBootSec=0h
OnUnitActiveSec=25h

[Install]
WantedBy=timers.target
