# ptc.timer
[Timer]
OnBootSec=0
OnUnitActiveSec=25h

[Install]
WantedBy=timers.target
